# supervisord.conf
[program:nginx]
command = /usr/sbin/nginx
autorestart = true
stdout_logfile = /data/logs/nginx-error.log
stderr_logfile = /data/logs/nginx-error.log
